# docker/development/etc/s6-overlay/s6-rc.d/scheduler-worker/run
#!/bin/sh

cd /var/www/html

if grep -qE '^SCHEDULER_ENABLED=false' .env 2>/dev/null; then
    echo "   INFO  Scheduler is disabled, sleeping."
    exec sleep infinity
fi

echo "   INFO  Scheduler is enabled, starting..."
exec php artisan schedule:work
